# deploy/supervisor.conf
[program:celery]
command=~/.env/bin/python3 python manage.py celery -l DEBUG
directory=/home/ubuntu/doc_disco
stdout_logfile=var/log/doc_disco/celery/celery-worker.log
stderr_logfile=var/log/doc_disco/celery/celery-worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998


[program:uwsgi]
user = uwsgi
command=/opt/bin/uwsgi --config=path to uwsgi.ini (/home/ubuntu/doc_disco/deploy/uwsgi.ini)
autostart=true
autorestart=true
stderr_logfile = var/log/doc_disco/uwsgi/err.log
stdout_logfile = var/log/doc_disco/out.log
stopsignal=INT
